# templates/default/storm-ui.conf.erb
# This file is auto-generated by chef
# please don't modify manually

description "storm ui process"

start on runlevel [2345]
stop on runlevel [016]
setuid <%= node["storm"]["user"] %>
respawn

exec <%= node["storm"]["symlink"] %>/bin/storm ui